# course_gather/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:frontend]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
directory=/app/webapp
command=npm start

[program:backend]
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
directory=/app
command=python3 manage.py runserver 0.0.0.0:8000
